# plugins/claude-tts-companion/launchd/com.terryli.claude-tts-companion.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.terryli.claude-tts-companion</string>
	<key>Program</key>
	<string>/Users/terryli/.local/bin/claude-tts-companion</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>NetworkState</key>
		<true/>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Nice</key>
	<integer>-5</integer>
	<key>ExitTimeOut</key>
	<integer>15</integer>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/terryli</string>
		<key>PATH</key>
		<string>/Users/terryli/.local/share/mise/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>SHERPA_ONNX_PATH</key>
		<string>/Users/terryli/fork-tools/sherpa-onnx/build-swift-macos/install</string>
		<key>KOKORO_MODEL_PATH</key>
		<string>/Users/terryli/.local/share/kokoro/models/kokoro-int8-multi-lang-v1_0</string>
		<key>TELEGRAM_BOT_TOKEN</key>
		<string></string>
		<key>TELEGRAM_CHAT_ID</key>
		<string></string>
		<key>MINIMAX_API_KEY</key>
		<string></string>
	</dict>
	<key>SoftResourceLimits</key>
	<dict>
		<key>MemoryLimit</key>
		<integer>8589934592</integer>
	</dict>
	<key>StandardErrorPath</key>
	<string>/Users/terryli/.local/state/launchd-logs/claude-tts-companion/stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/terryli/.local/state/launchd-logs/claude-tts-companion/stdout.log</string>
</dict>
</plist>
